# server-config/upstart-config/sk-app.conf
# Info
description "SK PHP Sample application"
author      "Mihkel Selgal"

# Events
start on filesystem and runlevel [2345]

# Automatically respawn
respawn
respawn limit 20 5

# Run the script!
# Note, in this example, if your PHP script returns
# the string "ERROR", the daemon will stop itself.
script
    [ $(exec /usr/bin/php -S localhost:8000 -t /apps/example/web) = 'ERROR' ] && ( stop; exit 1; )
end script